# salt/application/broker.supervisor.conf
[program:broker]
directory=/var/www/broker
user={{ pillar.deploy_user }}
command={% if pillar.env_name == 'vagrant' %}supervisor -i .git,client,node_modules,salt,scripts{% else %}node{% endif %} server
autostart=true
autorestart=true
{% if pillar.env_name != 'vagrant' %}environment=NODE_ENV=production{% endif %}
